# scripts/openvibe-tunnel.service
[Unit]
Description=OpenVibe SSH Tunnel to Huawei Cloud
Documentation=https://github.com/openvibe
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zcy
Group=acy2
Environment="AUTOSSH_GATETIME=0"

ExecStart=/usr/bin/autossh -M 0 -N -R 4096:localhost:4096 huawei \
    -o ServerAliveInterval=15 \
    -o ServerAliveCountMax=3 \
    -o StrictHostKeyChecking=no \
    -o BatchMode=yes \
    -o ExitOnForwardFailure=no

Restart=always
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=openvibe-tunnel

[Install]
WantedBy=multi-user.target
