# scripts/openvibe-tunnel.service
[Unit]
Description=OpenVibe SSH Tunnel to Huawei Cloud
Documentation=https://github.com/openvibe
After=network-online.target opencode.service
Wants=network-online.target
Requires=opencode.service

[Service]
Type=simple
User=zcy
Group=acy2

# SSH reverse tunnel command
# -N: No remote command
# -R: Remote port forwarding (Huawei:4096 -> Arch:4096)
ExecStart=/usr/bin/ssh -N -R 4096:localhost:4096 huawei \
    -o ServerAliveInterval=30 \
    -o ServerAliveCountMax=3 \
    -o ExitOnForwardFailure=yes \
    -o StrictHostKeyChecking=no \
    -o BatchMode=yes

# Restart on failure
Restart=always
RestartSec=10s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openvibe-tunnel

[Install]
WantedBy=multi-user.target
